# oncall/oncall-watcher.service
[Unit]
Description=aiNOC On-Call Watcher
Documentation=file:/home/mcp/aiNOC/README.md
After=network.target
Wants=vector.service

[Service]
Type=simple
User=mcp
WorkingDirectory=/home/mcp/aiNOC
ExecStart=/home/mcp/aiNOC/mcp/bin/python3 /home/mcp/aiNOC/oncall/watcher.py --service
Restart=on-failure
RestartSec=10

# Hardening
ProtectSystem=strict
ReadWritePaths=/home/mcp/aiNOC /var/log /tmp
ProtectHome=false

# Logging — integrates with journalctl -u oncall-watcher
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ainoc-watcher

# Environment — .env is loaded by Python's load_dotenv() at runtime.
# Do NOT use systemd EnvironmentFile here: systemd does not strip inline
# comments (KEY=val # comment), which corrupts values that python-dotenv
# handles correctly.
# EnvironmentFile=/home/mcp/aiNOC/.env

[Install]
WantedBy=multi-user.target
